# deploy/quizonline-env-fetch.service
[Unit]
Description=QuizOnline — fetch secrets from SSM Parameter Store
# Network has to be up (AWS API call) and clock has to be sane
# (SigV4 signatures rely on it). network-online.target is the
# canonical "we can talk to remote hosts" signal on systemd.
After=network-online.target
Wants=network-online.target

[Service]
# Type=oneshot + RemainAfterExit=yes makes this unit a *barrier*:
# Requires=quizonline-env-fetch.service on a downstream unit means
# the downstream waits for this to exit successfully, and the unit
# stays "active" thereafter so the dependency edge holds for the
# whole boot, not just the moment the script returns.
Type=oneshot
RemainAfterExit=yes

# Runs as root because it has to chmod/chown /run/quizonline/.env
# to django:www-data afterwards. The IAM identity used for the SSM
# call comes from the EC2 instance role via IMDS, not from any
# credentials on disk.
User=root

# Create /run/quizonline BEFORE ExecStart so the namespace setup
# below has something to mount. RuntimeDirectoryPreserve=yes keeps
# it alive after the oneshot returns, so downstream gunicorn /
# celery (which read /run/quizonline/.env via EnvironmentFile=)
# still see the file once env-fetch has finished. Default mode is
# 0755 root:root; the script chmods/chowns the .env file inside.
RuntimeDirectory=quizonline
RuntimeDirectoryPreserve=yes

# Hardening — /run/quizonline is the only writable path the script
# needs and RuntimeDirectory= above makes it implicitly RW even
# under ProtectSystem=strict, so no ReadWritePaths= is needed
# (and ReadWritePaths= would have failed: it requires the path
# to already exist at namespace setup, but on a fresh boot it
# doesn't until RuntimeDirectory creates it).
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
RestrictSUIDSGID=yes

# §3.10: this oneshot runs as root, so its script must NOT live in the
# django-writable tree — a compromised gunicorn (runs as django) could edit
# it and have it execute as root. It is installed root:root 0755 at
# /usr/local/sbin/quizonline-env-fetch.sh by ROOT from the committed git blob
# (see deploy.yml, same airtight pattern as the units). deploy/fetch-env-from-ssm.sh
# stays as the versioned source, never the execution target.
ExecStart=/usr/local/sbin/quizonline-env-fetch.sh

# Don't drown the journal if the script churns through retries.
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
