# deploy/quizonline-frontend-runtime-fetch.service
[Unit]
Description=QuizOnline — fetch frontend Sentry runtime config from SSM
# Same prerequisites as quizonline-env-fetch (network + clock for
# SigV4) plus nginx itself: the snippet is only useful once nginx is
# up to reload, and the ``nginx -t`` guard inside the script needs
# the nginx binary on PATH (always present once the package is
# installed, but ordering against nginx.service keeps the journal
# readable on a fresh boot).
After=network-online.target nginx.service
Wants=network-online.target

[Service]
# Type=oneshot mirrors quizonline-env-fetch — we don't want a
# long-lived process, just a one-shot snippet refresh + nginx
# reload. RemainAfterExit=yes keeps the unit in "active (exited)"
# state so downstream ``systemctl status`` shows a clear OK.
Type=oneshot
RemainAfterExit=yes

# Runs as root because (a) it writes /etc/nginx/snippets/, (b) the
# nginx config validation ``nginx -t`` needs to read root-owned
# certs/keys referenced from the main config, and (c) the
# ``systemctl reload nginx`` at the end needs root.
User=root

# Hardening — the script writes to /etc/nginx/snippets/ and needs
# /run accessible for ``nginx -t`` (probes /run/nginx.pid) and the
# subsequent ``systemctl reload nginx``. ProtectSystem=full leaves
# /usr /boot /efi /etc read-only but keeps /run writable — exactly
# what we need.
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/etc/nginx/snippets
ProtectHome=yes
NoNewPrivileges=yes
RestrictSUIDSGID=yes

# §3.10: this oneshot runs as root, so its script must NOT live in the
# django-writable tree. Installed root:root 0755 at
# /usr/local/sbin/quizonline-frontend-runtime-fetch.sh by ROOT from the committed
# git blob (see deploy.yml); deploy/fetch-frontend-runtime-from-ssm.sh is the
# versioned source, never the execution target.
ExecStart=/usr/local/sbin/quizonline-frontend-runtime-fetch.sh

TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
