# deploy/quizonline-gunicorn.service
[Unit]
Description=QuizOnline Django (gunicorn)
After=network.target quizonline-env-fetch.service
Requires=quizonline-env-fetch.service

[Service]
Type=simple
User=django
Group=www-data
WorkingDirectory=/var/www/django_websites/QuizOnline/quizonline-server

# Single source of truth: a tmpfs file written at boot by
# quizonline-env-fetch.service from AWS SSM Parameter Store.
#
# IMPORTANT: env-fetch is ``Type=oneshot RemainAfterExit=yes`` — the
# unit is "active" forever after the first run, so the Requires=
# barrier above is already satisfied and ``systemctl restart
# quizonline-gunicorn`` alone does NOT re-fetch from SSM. The full
# rotation flow is:
#
#   aws ssm put-parameter --overwrite ...
#   sudo systemctl restart quizonline-env-fetch   # re-runs ExecStart
#   sudo systemctl restart quizonline-gunicorn    # picks up new env
#
# Same applies to quizonline-celery / quizonline-celery-beat.
EnvironmentFile=/run/quizonline/.env
Environment="PYTHONUNBUFFERED=1"

ExecStart=/var/www/django_websites/QuizOnline/quizonline-server/.venv/bin/gunicorn \
    --workers 2 \
    --timeout 120 \
    --bind 127.0.0.1:8000 \
    --access-logfile /var/log/quizonline/gunicorn-access.log \
    --error-logfile /var/log/quizonline/gunicorn-error.log \
    config.wsgi:application

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
