# deploy/quizonline-celery-beat.service
[Unit]
Description=Celery Beat QuizOnline
After=network.target quizonline-env-fetch.service
Requires=quizonline-env-fetch.service

[Service]
User=django
Group=www-data
# Group www-data, not world-accessible (nginx still reads via the group).
UMask=0027
WorkingDirectory=/var/www/django_websites/QuizOnline/quizonline-server

# Secrets come from /run/quizonline/.env (tmpfs), populated at boot
# from AWS SSM Parameter Store. See quizonline-env-fetch.service.
EnvironmentFile=/run/quizonline/.env
Environment="PYTHONUNBUFFERED=1"

ExecStart=/var/www/django_websites/QuizOnline/quizonline-server/.venv/bin/celery \
    -A config beat \
    --pidfile=/var/www/django_websites/QuizOnline/quizonline-server/run/celerybeat.pid
    --schedule=/var/www/django_websites/QuizOnline/quizonline-server/celerybeat-schedule \
    --loglevel=info

Restart=always
RestartSec=5
LimitNOFILE=65536

PrivateTmp=true
[Install]
WantedBy=multi-user.target
