# deploy/quizonline-geoip-fetch.service
[Unit]
Description=QuizOnline — fetch MaxMind GeoLite2-City DB
# Needs network (MaxMind download) and the runtime env file produced by
# env-fetch: the script reads MAXMIND_* / GEOIP_PATH from /run/quizonline/.env,
# not from SSM directly (the box's default aws identity, certbot-route53, lacks
# ssm:GetParameter — see the script header).
After=network-online.target quizonline-env-fetch.service
Wants=network-online.target
Requires=quizonline-env-fetch.service

[Service]
Type=oneshot
# Runs as root to install the .mmdb owned django:www-data under GEOIP_PATH.
User=root

# ExecStart points straight at the git checkout (same pattern as
# quizonline-env-fetch.service) — no install into /usr/local/bin.
ExecStart=/var/www/django_websites/QuizOnline/deploy/fetch-geoip-from-ssm.sh

# A failed download is non-fatal: any existing DB stays in place and the
# connection-log geolocation simply degrades. Don't let a transient MaxMind
# hiccup mark the boot degraded.
SuccessExitStatus=0
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
