# deploy/quizonline-celery.service
[Unit]
Description=Celery Worker QuizOnline
After=network.target quizonline-env-fetch.service
Requires=quizonline-env-fetch.service

[Service]
Type=simple
User=django
Group=www-data
# Group www-data, not world-accessible (nginx still reads via the group).
UMask=0027
WorkingDirectory=/var/www/django_websites/QuizOnline/quizonline-server

# Secrets come from /run/quizonline/.env (tmpfs), populated at boot
# from AWS SSM Parameter Store. See quizonline-env-fetch.service.
EnvironmentFile=/run/quizonline/.env
Environment="PYTHONUNBUFFERED=1"

ExecStart=/var/www/django_websites/QuizOnline/quizonline-server/.venv/bin/celery \
    -A config worker \
    -n quizonline@%H \
    -Q quizonline \
    --loglevel=info \
    --concurrency=2

Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
